# crontab
22 10 * * * /Users/benwarren/Documents/GitHub/bosszp-selenium/boss_selenium.py
